# bind/rndc-reload.service
[Unit]
Description=rndc-reload

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
Type=simple
ExecStart=/usr/local/bin/rndc-reload
